# deploy/dbrain-watchdog.service
[Unit]
Description=d-brain watchdog (recovers a dead/hung Claude session)
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5
OnFailure=dbrain-notify@%n.service

[Service]
Type=simple
# KillMode=process: a watchdog restart must NOT kill the tmux brain it spawned
# via force_recover (the brain lives in this cgroup). Previously the default
# control-group + WatchdogSec SIGABRT loop killed the brain every minute.
KillMode=process
# Type=simple (not notify): the loop blocks during force_recover longer than a
# short WatchdogSec, so a systemd WatchdogSec would SIGABRT it mid-recovery and
# take the brain with it. Restart=always covers a crashed watchdog instead.
# A clean stop exits 143 (SIGTERM) — success, must not fire OnFailure=.
SuccessExitStatus=143 SIGTERM
# Separate slice so an OOM kill of the brain doesn't take the watchdog down.
Slice=dbrain-watchdog.slice
WorkingDirectory=%h/projects/dbrain
ExecStart=%h/.local/bin/uv run python -m d_brain.services.watchdog
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
Environment=CLAUDE_CONFIG_DIR=%h/.claude
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
